# qca/feeds/luci/themes/luci-theme-zx/root/etc/uci-defaults/luci-theme-zx
#!/bin/sh
uci batch <<-EOF
	set luci.themes.zx=/luci-static/zx
	set luci.main.mediaurlbase=/luci-static/zx
        commit luci
EOF
