# qca/feeds/luci/themes/luci-theme-slk/root/etc/uci-defaults/luci-theme-slk
#!/bin/sh
uci batch <<-EOF
	set luci.themes.slk=/luci-static/slk
	set luci.main.mediaurlbase=/luci-static/slk
        commit luci
EOF
