# qca/feeds/luci/applications/luci-app-serial-ipq/root/etc/init.d/ser2netd
#!/bin/sh /etc/rc.common
START=60

start()
{
		/etc/setusb 1 & 
		wifi up

}
