# qca/feeds/luci/themes/luci-theme-WitLinc/root/etc/uci-defaults/luci-theme-WitLinc
#!/bin/sh
uci batch <<-EOF
	set luci.themes.WitLinc=/luci-static/WitLinc
	set luci.main.mediaurlbase=/luci-static/WitLinc
        commit luci
EOF
