# qca/feeds/luci/themes/luci-theme-violet/root/etc/uci-defaults/30_luci-theme-violet
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Violet=/luci-static/violet
	set luci.main.mediaurlbase=/luci-static/violet
	commit luci
EOF
exit 0
